# docker/samba-ad-dc/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/samba/supervisord.log

[program:samba]
command=/usr/sbin/samba --foreground --no-process-group
autostart=true
autorestart=true
priority=10
stdout_logfile=/data/logs/samba/samba-stdout.log
stderr_logfile=/data/logs/samba/samba-stderr.log

[program:sudo-schema]
command=/usr/local/lib/setup-sudo-schema.sh
autostart=true
autorestart=false
startsecs=0
priority=20
stdout_logfile=/data/logs/samba/sudo-schema-stdout.log
stderr_logfile=/data/logs/samba/sudo-schema-stderr.log
