# docker/all-in-one/supervisor/sudo-schema.conf
[program:sudo-schema]
command=/usr/local/lib/setup-sudo-schema.sh
autostart=true
autorestart=false
startsecs=0
priority=15
stdout_logfile=/data/logs/samba/sudo-schema-stdout.log
stderr_logfile=/data/logs/samba/sudo-schema-stderr.log
